# setup/setup_files/xhost_enable_local_user.service
[Unit]
Description=Allow Docker containers to access XWayland after Wayland login
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
Environment=DISPLAY=:0
ExecStart=/usr/bin/xhost +local:
RemainAfterExit=yes

[Install]
WantedBy=default.target
