# setup/setup_files/xhost_enable_local.service
[Unit]
Description=Allow Docker containers to access the host XWayland display
After=graphical.target display-manager.service
Wants=graphical.target

[Service]
Type=oneshot
User=pib
Environment=DISPLAY=:0
ExecStart=/usr/bin/xhost +local:
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
